# cron.d
* * * * * /scripts/update_trends.sh >> /var/log/cron.log 2>&1


